# src/dist/vsphere2metrics.service
[Unit]
Description=vsphere2metrics
After=network.target

[Service]
Type=simple
RuntimeDirectory=/opt/vsphere2metrics
Environment=JAVA_HOME=/usr/java/latest
ExecStart=/opt/vsphere2metrics/vsphere2metrics
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=60
StandardOutput=journal
StandardError=journal
LimitNOFILE=infinity
LimitMEMLOCK=infinity
LimitNPROC=infinity
LimitAS=infinity

[Install]
WantedBy=multi-user.target
